# src/random_quotes.service
[Unit]
Description=Writes a random quote in /etc/issue

[Service]
ExecStart=/usr/bin/random_quote

[Install]
WantedBy=multi-user.target
